# systemd-stuff/test3.service
[Unit]
Description=asdf
BindsTo=test1.service test2.service

[Service]
ExecStart=/path/to/executable3
Restart=on-failure
