# systemd-stuff/test5.service
[Unit]
Description=test5 service for script5.py

[Service]
Type=oneshot
ExecStart=/path/to/script5.py

[Install]
WantedBy=multi-user.target
