# systemd-stuff/test2.service
[Unit]
Description=asdf
BindsTo=test1.service test3.service

[Service]
ExecStart=/path/to/executable2
Restart=on-failure
